# addons/systemd/hmonitor.service
[Unit]
Description=Analyze nagios hoststate
Before=nagios4.service

[Service]
Type=forking
User=nagios
Group=nagios
ExecStart=/opt/Eye/scripts/hmonitor.pl --start
ExecStop=/opt/Eye/scripts/hmonitor.pl --stop
RestartSec=30
Restart=on-failure
PIDFile=/run/nagios4/hmonitor.pid

[Install]
WantedBy=multi-user.target
